# server/systemd/timbranfc-server.service
[Unit]
Description=TimbraNFC Server di Sede
After=network.target mariadb.service mysql.service
Wants=mariadb.service

[Service]
User=timbranfc
Group=timbranfc
WorkingDirectory=/opt/timbranfc
EnvironmentFile=-/opt/timbranfc/.env
ExecStart=/opt/timbranfc/.venv/bin/python -m uvicorn server.app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
