# systemd/dashboard.service
[Unit]
Description=Dashboard web presenze
After=network.target timbratrice.service

[Service]
User=pi
WorkingDirectory=/home/pi/TimbraNFC
Environment=ADMIN_PASSWORD=cambia-in-produzione
Environment=DASHBOARD_SECRET=chiave-segreta-random
ExecStart=/usr/bin/python3 dashboard/app.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
