# standalone/systemd/timbranfc-server.service
# Template di riferimento — NON installare manualmente.
# I servizi reali sono generati da install-raspberry.sh o fix-services.sh
# con User/WorkingDirectory/ExecStart per l'utente corrente.

[Unit]
Description=TimbraNFC — Server e dashboard web (LAN)
After=network.target
Before=timbranfc-kiosk.service

[Service]
Type=simple
User=APP_USER
Group=APP_USER
WorkingDirectory=APP_DIR
EnvironmentFile=-APP_DIR/.env
Environment=STANDALONE=1
Environment=TIMBRANFC_DATA=APP_DIR/data
Environment=SERVER_HOST=0.0.0.0
Environment=SERVER_PORT=8080
ExecStart=APP_DIR/.venv/bin/python APP_DIR/standalone/run_server.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
