# standalone/autostart/timbranfc-kiosk.desktop
[Desktop Entry]
Type=Application
Name=TimbraNFC Kiosk
Comment=Timbratrice presenze NFC — avvio automatico (fallback)
Exec=@APP_DIR@/standalone/launch_kiosk.sh
Terminal=false
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
X-GNOME-Autostart-Delay=25
